# etc/systemd/system/gw-watchdog.timer
[Unit]
Description=IoTmaxx device watchdog timer

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target
